# deploy/prismd.service
[Unit]
Description=PrismDB multi-model database server
Documentation=https://github.com/HafizMMoaz/prism-db
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# A transient, unprivileged service account is allocated per start; combined with
# StateDirectory this needs no useradd in the package's install scripts.
DynamicUser=true

# systemd creates and owns /var/lib/prismdb for the service. The config file
# (installed at /etc/prismdb/prismd.conf) points PRISM_DATA_DIR at it and sets
# the bind address, log level, and optional TLS — edit it, then restart.
StateDirectory=prismdb
EnvironmentFile=/etc/prismdb/prismd.conf

# All configuration (bind address, data dir, TLS) comes from the environment
# file above; prismd reads PRISM_BIND / PRISM_DATA_DIR / PRISM_TLS_CERT /
# PRISM_TLS_KEY / RUST_LOG.
ExecStart=/usr/bin/prismd run

Restart=on-failure
RestartSec=2

# Hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
